# deploy/jobmarket-ingest.service
[Unit]
Description=Job Market Intelligence — ingest-only run (collect + normalize + dedupe)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/docker exec jobmarket-web python -m src.orchestrator --mode ingest-only
StandardOutput=journal
StandardError=journal
# Ingest-only should complete well within an hour
TimeoutStartSec=3600
