# deploy/jobmarket-ingest.timer
[Unit]
Description=Run Job Market Intelligence ingestion every 12 hours
Requires=jobmarket-ingest.service

[Timer]
# First run 5 minutes after the timer unit starts (e.g. after reboot)
OnBootSec=5min
# Then every 12 hours after the previous run finishes
OnUnitActiveSec=12h
# Run at a fixed wall-clock time so both VPS instances stay in sync (optional)
# OnCalendar=*-*-* 06,18:00:00
# Randomise start by up to 10 minutes to avoid hammering APIs on the dot
RandomizedDelaySec=600
Unit=jobmarket-ingest.service

[Install]
WantedBy=timers.target
